# supervisor/conf.d/nginx.conf.template
[program:nginx]
directory=${V2RAY_PATH}
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stopsignal=INT
