# supervisor/conf.d/v2fly.conf.template
[program:v2fly]
directory=${V2RAY_PATH}
command=v2ray
autostart=true
autorestart=true
stopsignal=INT